# crontab
01 03 * * * /home/tom/backbot/backup.sh
